# services/cas-rf-u2021xa-visa.service
[Unit]
Description=Sirius RF U2021XA Power Meter - Visa Interface
Wants=network-online.target
After=network-online.target

[Service]
Restart=always
RestartSec=10
Type=simple

WorkingDirectory=/opt/cas-rf-u2021xa/scripts
ExecStart=/bin/bash -c "\
     set -x;\
     set -e;\
     source /root/miniconda3/etc/profile.d/conda.sh;\
     conda activate base;\
     python /opt/cas-rf-u2021xa/scripts/run.py\
         --resource '${RESOURCE}'\
"

StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
